# docker/supervisord.conf
; Runs the four services plus nginx in one container. Container Apps exposes a
; single port, so nginx (8080) is the only thing the ingress talks to.
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
user=root

[program:django]
command=bash /app/docker/start-django.sh
directory=/app/Backend
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nextjs]
command=node server.js
directory=/app/frontend
environment=PORT="3000",HOSTNAME="127.0.0.1"
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:face]
command=uvicorn face_embedding:app --host 127.0.0.1 --port 8001
directory=/app/AI
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; The agent logs into Django on startup and exits if that fails, so it may lose
; the race on a cold boot. autorestart lets it retry until Django is listening.
[program:agent]
command=uvicorn server:fastapi_app --host 127.0.0.1 --port 8090
directory=/app/AI/Agent
autorestart=true
startretries=20
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;" -c /app/docker/nginx.conf
autorestart=true
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
